# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
; supervisord 自己的 pid 写到 /tmp，不占 /app/data
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; ── 监控进程 ─────────────────────────────────────────────────────
[program:monitor]
command=python -u monitor.py
directory=/app
autostart=true
autorestart=true
; 崩溃后最多重启 10 次；若持续崩溃由 Docker healthcheck 兜底重启容器
startretries=10
; 优雅停止：先等 10 秒让 monitor 完成当前轮次
stopwaitsecs=10
; monitor 内部用 Python logging 写 data/monitor.log，不重复 capture
redirect_stderr=true
stdout_logfile=/app/logs/monitor_stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3

; ── Web 面板 ─────────────────────────────────────────────────────
[program:web]
; 使用 Gunicorn 生产 WSGI 服务器（替代 Flask 内置的 Werkzeug 开发服务器）
; --workers=1   : SQLite 单进程，避免多进程写锁冲突
; --threads=8   : 每个 worker 8 条线程，支持多路 SSE 长连接并发
; --timeout=0   : 禁用 worker 超时，防止 SSE 长连接被强杀（默认 30s 会断连）
; web:app       : web.py 中的 Flask WSGI application 对象
command=gunicorn --workers=1 --threads=8 --timeout=0 --bind=0.0.0.0:8088 web:app
directory=/app
autostart=true
autorestart=true
startretries=3
stopwaitsecs=10
redirect_stderr=true
stdout_logfile=/app/logs/web.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
